# docs/85-openjvs.rules
ACTION=="add", DRIVERS=="usbhid", ENV{ID_BUS}=="usb", RUN+="/bin/systemctl restart openjvs"
